# deploy/roasti-backup.service
[Unit]
Description=Roasti Database Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=roasti
EnvironmentFile=/var/lib/roasti/backup.env
ExecStart=/usr/local/bin/roasti-backup.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=roasti-backup

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/roasti /tmp
